# shared/com.jarvis.xhs-feedback.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jarvis.xhs-feedback</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jarvis/xiaohongshu-bot/shared/feedback.py</string>
    </array>

    <!-- 每天 02:00 运行（两个 MCP server 此时必然空闲）-->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/jarvis/xiaohongshu-bot/shared/feedback-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jarvis/xiaohongshu-bot/shared/feedback-launchd.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin</string>
        <key>HOME</key>
        <string>/Users/jarvis</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/jarvis/xiaohongshu-bot/shared</string>
</dict>
</plist>
