# launchd/com.jarvis.xhs-curate-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jarvis.xhs-curate-weekly</string>

    <!-- 每周日 07:00 从权威期权教育账号精选素材 → state/curated/YYYY-WW-options.json -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/jarvis/xiaohongshu-mcp &amp;&amp; /usr/bin/python3 pipeline/track_a/curator.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/jarvis/xiaohongshu-mcp/logs/curate-weekly.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jarvis/xiaohongshu-mcp/logs/curate-weekly.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/Users/jarvis/.local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/jarvis</string>
    </dict>
</dict>
</plist>
